# scripts/com.robotrader.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.robotrader.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/oliver/robo_trader/scripts/watchdog.sh</string>
        <string>5</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/oliver/robo_trader</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/oliver/robo_trader/watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/oliver/robo_trader/watchdog.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
